# debian/udoo-web-conf.upstart

#
# udoo-wcp - start script
#

description "UDOO Neo Web Configurator"
author "UDOO Team <social@udoo.org>"

env USER=root
env PATH=/sbin:/bin:/usr/bin:/usr/sbin

start on local-filesystems
stop on [06]

respawn
respawn limit 4 20

script
  chdir /opt/udoo-wcp/public
  exec start-stop-daemon --chdir /opt/udoo-wcp/public \
    --chuid root \
    --start --make-pidfile \
    --pidfile /var/run/udoo-wcp.pid \
    --exec php -- -S 0.0.0.0:$(cat /etc/udoo-wcp/port)
end script
