# debian/udoo-wcp.upstart

#
# udoo-wcp - start script
#

description "UDOO Neo Web Configurator"
author "UDOO Team <social@udoo.org>"

start on local-filesystems
stop on [06]

respawn
respawn limit 4 20

script
  chdir /opt/udoo-wcp/public
  exec start-stop-daemon --chdir /opt/udoo-wcp/public \
    --chuid root \
    --start \
    --exec /usr/bin/php -- -S 0.0.0.0:$(cat /etc/udoo-wcp/port)
end script
